# crontab
30 12 * * * python /usr/src/app/crwal.py